# debian/postinst
#!/bin/bash

update-rc.d cron-ms defaults
chmod 755 /etc/init.d/cron-ms

